# auditd/hardened-homelab.rules
## Optimized auditd rules for homelab - ERROR FREE
## Tested on Ubuntu 24.04 with kernel 6.x

# ============================================================
# CONFIGURATION
# ============================================================

-D
-b 32768
-f 1
--backlog_wait_time 60000
-r 0

# ============================================================
# CRITICAL FILE MONITORING
# ============================================================

## Identity files
-w /etc/passwd -p wa -k identity
-w /etc/shadow -p wa -k identity
-w /etc/group -p wa -k identity
-w /etc/gshadow -p wa -k identity
-w /etc/security/opasswd -p wa -k identity

## SSH configuration
-w /etc/ssh/sshd_config -p wa -k sshd_config
-w /etc/ssh/ssh_config -p wa -k ssh_config
-w /root/.ssh -p wa -k root_ssh

## Sudo configuration
-w /etc/sudoers -p wa -k sudoers
-w /etc/sudoers.d -p wa -k sudoers

## PAM configuration
-w /etc/pam.d -p wa -k pam_config

## Login tracking
-w /var/log/lastlog -p wa -k login
-w /var/log/faillog -p wa -k login
-w /var/run/utmp -p wa -k session
-w /var/log/wtmp -p wa -k session
-w /var/log/btmp -p wa -k session

# ============================================================
# NETWORK & FIREWALL
# ============================================================

## nftables
-w /etc/nftables.conf -p wa -k firewall

## Network configuration
-w /etc/network -p wa -k network_config
-w /etc/netplan -p wa -k network_config
-w /etc/systemd/network -p wa -k network_config

## Hosts and DNS
-w /etc/hosts -p wa -k network_config
-w /etc/hostname -p wa -k network_config
-w /etc/resolv.conf -p wa -k network_config

## Tailscale
-w /etc/default/tailscaled -p wa -k tailscale
-w /var/lib/tailscale -p wa -k tailscale

# ============================================================
# CONTAINER INFRASTRUCTURE
# ============================================================

## Incus/LXD
-w /var/lib/incus -p wa -k containers
-w /etc/subuid -p wa -k containers
-w /etc/subgid -p wa -k containers

## Docker
-w /etc/docker -p wa -k containers

## Podman
-w /etc/containers -p wa -k containers

# ============================================================
# SECURITY FRAMEWORKS
# ============================================================

## AppArmor
-w /etc/apparmor.d -p wa -k apparmor
-w /etc/apparmor -p wa -k apparmor
-w /sys/kernel/security/apparmor -p wa -k apparmor

## CrowdSec
-w /etc/crowdsec -p wa -k crowdsec

## AIDE
-w /etc/aide -p wa -k aide
-w /var/lib/aide -p wa -k aide

## Audit self-monitoring
-w /etc/audit -p wa -k auditd_config
-w /sbin/auditctl -p x -k auditd_tools
-w /sbin/auditd -p x -k auditd_tools

# ============================================================
# SYSTEM CONFIGURATION
# ============================================================

## Boot configuration
-w /boot -p wa -k boot
-w /etc/grub.d -p wa -k boot
-w /etc/default/grub -p wa -k boot

## Systemd
-w /etc/systemd/system -p wa -k systemd
-w /lib/systemd/system -p wa -k systemd

## Cron jobs
-w /etc/cron.allow -p wa -k cron
-w /etc/cron.deny -p wa -k cron
-w /etc/cron.d -p wa -k cron
-w /etc/cron.daily -p wa -k cron
-w /etc/cron.hourly -p wa -k cron
-w /etc/cron.monthly -p wa -k cron
-w /etc/cron.weekly -p wa -k cron
-w /etc/crontab -p wa -k cron
-w /var/spool/cron -p wa -k cron

## Init scripts (Ubuntu/Debian)
-w /etc/init.d -p wa -k init

## Kernel parameters (sysctl replacement - no syscall rule)
-w /etc/sysctl.conf -p wa -k sysctl
-w /etc/sysctl.d -p wa -k sysctl

## Critical binaries
-w /usr/bin/sudo -p wa -k critical_binary
-w /usr/sbin/sshd -p wa -k critical_binary

# ============================================================
# SYSTEM CALL MONITORING
# ============================================================

## Kernel modules (rootkit detection)
-w /sbin/insmod -p x -k modules
-w /sbin/rmmod -p x -k modules
-w /sbin/modprobe -p x -k modules
-a always,exit -F arch=b64 -S init_module,delete_module,finit_module -k modules

## Privilege escalation (UID >= 1000 = interactive users only)
-a always,exit -F arch=b64 -S execve -C uid!=euid -F euid=0 -F auid>=1000 -k privesc
-a always,exit -F arch=b64 -S execve -C gid!=egid -F egid=0 -F auid>=1000 -k privesc

## Sudo usage (interactive users only)
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/sudo -F auid>=1000 -k sudo

## Su usage
-a always,exit -F arch=b64 -S execve -F path=/bin/su -F auid>=1000 -k su
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/su -F auid>=1000 -k su

## Capability changes
-a always,exit -F arch=b64 -S capset -k capabilities

## Container operations (escape detection)
-a always,exit -F arch=b64 -S unshare,setns -k container_escape
-a always,exit -F arch=b64 -S mount,umount2 -F auid>=1000 -k mounts
-a always,exit -F arch=b64 -S pivot_root -k pivot_root
-a always,exit -F arch=b64 -S chroot -k chroot

## Time changes
-a always,exit -F arch=b64 -S adjtimex,settimeofday,clock_settime -k time_change

## Process tracing (debugger detection)
-a always,exit -F arch=b64 -S ptrace -F auid>=1000 -k ptrace

## Process memory operations (injection detection)
-a always,exit -F arch=b64 -S process_vm_writev,process_vm_readv -k process_memory

## Extended attributes (AppArmor bypass attempts)
-a always,exit -F arch=b64 -S setxattr,lsetxattr,fsetxattr -F auid>=1000 -k xattr_modification
-a always,exit -F arch=b64 -S removexattr,lremovexattr,fremovexattr -F auid>=1000 -k xattr_modification

## Execution from tmp (malware staging detection)
-w /tmp -p x -k tmp_exec
-w /var/tmp -p x -k tmp_exec
-w /dev/shm -p x -k tmp_exec

# ============================================================
# HOMELAB SPECIFIC
# ============================================================

## Docker Compose
-w /usr/bin/docker-compose -p x -k docker_compose
-w /usr/local/bin/docker-compose -p x -k docker_compose

## Incus commands
-w /usr/bin/incus -p x -k incus_command

## gVisor
-w /usr/local/bin/runsc -p x -k gvisor

# ============================================================
# IMMUTABLE (Enable after 1-2 weeks of testing)
# ============================================================

## Uncomment to make audit configuration immutable:
# -e 2